# systemd/system/wslogd.service.in
[Unit]
Description=Weather Station Logger

[Service]
User=wslog
ExecStart=@prefix@/bin/wslogd
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
