# packages/archlinux/wslogd.service
[Unit]
Description=Weather station logger

[Service]
User=wslog
ExecStart=/usr/bin/wslogd
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
